# systemd/system/pdns.service.d/hardening.conf
[Service]
# incompatible with LuaJIT
MemoryDenyWriteExecute=true

PrivateIPC=true
ProtectProc=invisible
ProtectSystem=strict
RemoveIPC=true
