# systemd/system/pdns.service.d/hardening.conf
[Service]
PrivateIPC=true
ProtectProc=invisible
ProtectSystem=strict
RemoveIPC=true
