# systemd/system/pdns.service.d/local.conf
[Service]
LimitNOFILE=2048
PrivateIPC=true
ProtectProc=invisible
ProtectSystem=strict
RemoveIPC=true
Restart=always
