# systemd/system/geoipupdate.timer.d/daily.conf
[Unit]
Description=Daily update of GeoIP2 and GeoIP Legacy databases

[Timer]
OnCalendar=
OnCalendar=daily
