# systemd/system/ipinfo-update.timer
[Unit]
Description=Daily IPinfo data update

[Timer]
OnCalendar=daily
Persistent=yes

[Install]
WantedBy=timers.target
