# systemd/system/bird.service.d/override.conf
[Unit]
After=pdns.service dnsdist.service
