# systemd/system/bird.service.d/override.conf
[Unit]
After=zerotier-one.service pdns.service dnsdist.service
