# systemd/system/geoipupdate.service.d/hardening.conf
[Service]
PrivateIPC=true
ProcSubset=pid
User=geoipupdate
