# systemd/system/geoipupdate.service.d/local.conf
[Service]
CPUSchedulingPolicy=batch
PrivateIPC=true
ProcSubset=pid
User=geoipupdate
