# systemd/system/geoipupdate.service.d/local.conf
[Service]
CPUSchedulingPolicy=batch
ExecStartPost=+/usr/bin/pdns_control reload
PrivateIPC=true
ProcSubset=pid
User=geoipupdate
