# systemd/system/geoipupdate.service.d/override.conf
[Service]
CPUSchedulingPolicy=batch
ExecStartPost=+/usr/bin/pdns_control reload
IOSchedulingPriority=7
Nice=19
PrivateIPC=true
ProcSubset=pid
User=geoipupdate
